# supervisor/mariadb.conf
[program:mariadb]
command=/usr/bin/pidproxy /usr/sbin/mysqld_safe
