# supervisor/mariadb.conf
[program:mariadb]
command=/usr/bin/pidproxy /var/run/mysqld.pid /usr/bin/mysqld_safe
stdout_logfile=/var/log/mysql/supervisor.log
stderr_logfile=/var/log/mysql/supervisor-error.log
